# ttrss/rootfs/etc/services.d/caddy/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Home Assistant Add-on: TT-RSS
# Runs Caddy
# ==============================================================================

bashio::log.info 'Starting Caddy...'

# Run Caddy
/usr/bin/caddy run --config /etc/caddy/Caddyfile
